# deployment/chef/cookbooks/ss-sync/templates/default/ss-sync.conf.erb
#!upstart
description "ss-sync service"
author      "SpeedySeat"

env NODE_PATH=/usr/local/lib/node_modules

respawn
start on runlevel [2345]
stop on runlevel [016]

exec /usr/local/bin/coffee /home/www-server/onthelist/sync/app.coffee
