# deployment/chef/cookbooks/authorize-device/templates/default/authorize-device.conf.erb
#!upstart
description "authorize-device service"
author      "SpeedySeat"

respawn
start on runlevel [2345]
stop on runlevel [016]

script
    export NODE_PATH=/usr/local/lib/node_modules

    exec /usr/local/bin/coffee /home/www-server/onthelist/auth/device/app.coffee 2>&1 >> /var/log/authorize-device.log
end script
