# deployment/chef/cookbooks/authorize-device/templates/default/authorize-device.conf.erb
exec /usr/local/bin/coffee /home/www-server/onthelist/auth/device/app.coffee

respawn
start on runlevel [2345]
stop on runlevel [016]
