# deployment/chef/cookbooks/ss-notifier/templates/default/ss-notifier.conf.erb
exec /usr/local/bin/coffee /home/www-server/onthelist/notify/app.coffee

respawn
start on runlevel [2345]
stop on runlevel [016]
