# notify/conf/ss-notifier.conf
exec /usr/local/bin/coffee /home/www-server/onthelist/notify/app.coffee

respawn
start on started network-services
stop on stopping network-services


